# deploy/nba-dashboard.service
[Unit]
Description=NBA Betting Dashboard (Streamlit)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/nbamodels
Environment="PATH=/root/nbamodels/.venv/bin"
ExecStart=/root/nbamodels/.venv/bin/streamlit run analytics_dashboard.py --server.port 8501 --server.address 0.0.0.0
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
